# conky/.config/autostart/conky.desktop
[Desktop Entry]
Exec=/usr/bin/conky
Name=conky
Type=Application
