# redshift/.config/autostart/redshift-gtk.desktop
[Desktop Entry]
Exec=/usr/bin/redshift-gtk
Name=redshift-gtk
Type=Application
